# deploy/com.jscraft.deploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jscraft.deploy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jongsoobae/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>--env-file=/opt/jscraft/jscraft-infra/deploy/.env</string>
        <string>/opt/jscraft/jscraft-infra/deploy/index.js</string>
    </array>
    <key>UserName</key>
    <string>jongsoobae</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jongsoobae/.nvm/versions/node/v24.14.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/opt/jscraft/jscraft-infra/deploy</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/opt/jscraft/logs/deploy.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/jscraft/logs/deploy.error.log</string>
</dict>
</plist>
